# xibo/xibo-server/vendor/xibosignage/xibo-xmr/bin/upstart/xmr.conf
description     "XMR Message Relay"

start on runlevel [3]
stop on runlevel [!2345]

respawn

script
    su -c "/usr/bin/php /opt/xmr/xmr.phar" -s /bin/sh nobody
end script